# supervisor-tabungan.conf
[program:tabungan-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tabungan/artisan queue:listen --tries=3 --delay=10 --sleep=5 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/tabungan/storage/logs/queue-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:tabungan-reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tabungan/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/tabungan/storage/logs/reverb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
